# init/systemd/zapret-daemon.service
[Unit]
Description=Zapret Daemon Service
Documentation=https://github.com/zapret/zapret-daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/zapret-daemon serve --config /etc/zapret/config.yaml
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run/zapret

# Capabilities needed for iptables manipulation
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# User and group
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zapret-daemon

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
